# tmpreaper/etc/cron.d/tmpreaper
# run tmpreaper 5h30 each days
30 05 * * * <?php echo getenv('APPNAME_HOME') ?>/usr/sbin/tmpreaper > /dev/null

