# postgresql-8.3/etc/init.d/postgresql-8.3
#!/bin/sh -e

### BEGIN INIT INFO
# Provides:		postgresql postgresql-8.4
# Required-Start:	$local_fs $remote_fs $network $time
# Required-Stop:	$local_fs $remote_fs $network $time
# Should-Start:		$syslog
# Should-Stop:		$syslog
# Default-Start:	2 3 4 5
# Default-Stop:		0 1 6
# Short-Description:	PostgreSQL 8.4 RDBMS server
### END INIT INFO

# Setting environment variables for the postmaster here does not work; please
# set them in /etc/postgresql/8.4/<cluster>/environment instead.

. <?php echo getenv('APPNAME_HOME') ?>/etc/profile

[ -r <?php echo getenv('APPNAME_HOME') ?>/usr/share/postgresql-common/init.d-functions ] || exit 0

. <?php echo getenv('APPNAME_HOME') ?>/usr/share/postgresql-common/init.d-functions

VERSION=${PGVERSION}


check_status() {
        pg_ctl status -w
        RESULT=`echo $?`
        if [ "$RESULT" = "0" ]; then
            log_action_msg "PostgreSQL is running (listening on <?php echo getenv('APPNAME_PG_LISTEN_INTERFACE')?>:<?php echo getenv('APPNAME_PG_LISTEN_PORT')?>)"
            exit 0
        else
            log_action_msg "PostgreSQL is stopped"
            exit 1
        fi
}


case "$1" in
    start)
        pg_ctl start -w
        check_status
        ;;
    stop)
        pg_ctl stop -w
        ;;
    restart)
        pg_ctl restart -w
        ;;
    force-reload | reload)
        pg_ctl reload -w
        ;;
    status)
        check_status
	;;
    *)
        echo "Usage: $0 {start|stop|restart|reload|force-reload|status}"
        exit 1
        ;;
esac

exit 0

