# php/etc/profile.d/php
export PHPRC="<?php echo getenv('APPNAME_HOME') ?>/etc/php5/php.ini"
