# core/etc/profile.d/10-core
PATH="<?php echo getenv('APPNAME_HOME') ?>/usr/sbin:<?php echo getenv('APPNAME_HOME') ?>/usr/bin:$PATH"

export PATH